# _experiments/hyperpod_imds_proxy/imds-proxy-server.service
[Unit]

Description=IMDS Proxy service

After=sagemaker-cluster-after-network-online.service

Before=imds-proxy-iptables.service
Before=slurmd.service
Before=sshd.service

[Service]
Type=simple
WorkingDirectory=/usr/local/bin
ExecStart=/usr/bin/python3 -u /usr/local/bin/imds-proxy-server.py

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
