# _experiments/hyperpod_imds_proxy/imds-proxy-iptables.service
[Unit]

Description=IMDS Proxy iptables configuration service

After=network-online.target
Wants=network-online.target

After=sagemaker-role-proxy-agent.service
Wants=sagemaker-role-proxy-agent.service

After=sagemaker-cluster-after-network-online.service
Wants=sagemaker-cluster-after-network-online.service

# Before=slurmd.service
# Before=sshd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/imds-proxy-iptables.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target