# _experiments/hyperpod_slurm_custom_health_monitoring/custom-health-monitor.service
[Unit]
Description=HyperPod Slurm Custom Health Monitor
After=network.target slurmd.service
Wants=slurmd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/custom-health-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=custom-health-monitor

# Security settings
User=root
Group=root

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
