# _experiments/hyperpod_slurm_shutdown_action/shutdown-action.service
[Unit]

Description=Shutdown action test service
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStop=/usr/bin/python3 -u /usr/local/bin/shutdown-action.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
